# config/supervisor/solvecore.conf
[program:solvecore]
command=/home/solvecore/solvecore/scripts/runserver.bash
user=solvecore
